# scripts/fast-track-deploy.timer
[Unit]
Description=Fast Track Auto Deploy Timer
After=network-online.target

[Timer]
# Run every 5 minutes
OnBootSec=1min
OnUnitActiveSec=5min

# Add some randomness to prevent thundering herd
RandomizedDelaySec=30

[Install]
WantedBy=timers.target
