# scripts/fast-track-deploy.service
[Unit]
Description=Fast Track Auto Deploy
After=network-online.target docker.service
Wants=network-online.target

[Service]
Type=oneshot
User=theflyingdutchman
WorkingDirectory=/opt/fast-track
ExecStart=/opt/fast-track/scripts/deploy.sh
Environment=REPO_DIR=/opt/fast-track
Environment=BRANCH=main

# Logging
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
